# trqauthd.service
[Unit]
Description=Torque Authorization Daemon

[Service]
Type=forking
ExecStart=/usr/sbin/trqauthd

[Install]
WantedBy=multi-user.target 
